# SOC-Lab/4-detection-engineering/persistence_cron_detection.spl
Persistence Detection (Cron Job)
🎯 MITRE
T1053
🔍 Logic

Detect cron modifications

💻 SPL
index=* "crontab"
| stats count by user, command
📊 Example
crontab -e
⚠️ False Positives
system maintenance jobs
⚙️ Tuning
alert only for suspicious commands
🔥 Severity-HIGH
